# inFight_pi_snapshot/configs/inflight-client.service
[Unit]
Description=inFlight Client service
After=inflight-server.service
Requires=inflight-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/inflight/frontend/client
ExecStart=/usr/bin/npm start --prefix /opt/inflight/frontend/client
#Restart=always # or always, on-abort, etc
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
